# 8006/Assignment 3/ips.service
[Unit]
Description=IPS
After=sshd.service

[Service]
ExecStart=/usr/local/etc/runner.sh

[Install]
WantedBy=user.target